# smsgw-watchdog.cron
*/5 * * * * root /usr/local/bin/smsgw-watchdog.sh >> /var/log/smsgw-watchdog.log 2>&1
