# smsgw-watchdog.cron
*/5 * * * * /usr/local/bin/smsgw-watchdog.sh
